# ks_neo.conf

description	"kickstarter monitor server"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
respawn limit 10 5

exec /home/udooer/ks_x86_neo.py

